# etc/cron.d/bootmail
@reboot [ -x /usr/sbin/bootmail ] && /usr/sbin/bootmail
